# susu-brain-backend.service
[Unit]
Description=Susu Brain Backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/html
EnvironmentFile=/etc/wa-agent.env
ExecStart=/usr/bin/python3 /var/www/html/susu_brain_backend.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
